# systemd/luminos-ai.service
# luminos-ai.service — Luminos OS central IPC daemon (Go, Phase 1)
# Replaces: old Python daemon/main.py service (Decision 13 — Go rewrite)
# Install: sudo cp this /etc/systemd/system/ && sudo systemctl daemon-reload
# Enable:  sudo systemctl enable --now luminos-ai
# [CHANGE: claude-code | 2026-04-20] Phase 1 Go foundation — replaced Python service.

[Unit]
Description=Luminos OS — Central AI Router Daemon (Go)
Documentation=https://github.com/luminos-os/luminos
After=network.target
Wants=luminos-power.service luminos-sentinel.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/luminos-ai
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luminos-ai
# RuntimeDirectory creates /run/luminos at service start so sub-daemons can bind sockets.
RuntimeDirectory=luminos
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
