# systemd/luminos-sentinel.service
# luminos-sentinel.service — Luminos OS security process scanner (Go, Phase 1)
# Replaces: old Python sentinel_daemon service (Decision 13 — Go rewrite)
# Install: sudo cp this /etc/systemd/system/ && sudo systemctl daemon-reload
# Enable:  sudo systemctl enable --now luminos-sentinel
# [CHANGE: claude-code | 2026-04-20] Phase 1 Go foundation — replaced Python service.

[Unit]
Description=Luminos OS — Rule-Based Security Process Scanner (Go)
Documentation=https://github.com/luminos-os/luminos
After=network.target luminos-ai.service

[Service]
Type=simple
ExecStart=/usr/local/bin/luminos-sentinel
Restart=always
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luminos-sentinel
# CAP_SYS_PTRACE lets the sentinel read /proc/[pid]/exe for non-owned processes.
AmbientCapabilities=CAP_SYS_PTRACE
NoNewPrivileges=false
RuntimeDirectory=luminos
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
