# systemd/luminos-power.service
[Unit]
Description=Luminos OS — Automatic Power & Thermal Manager
Documentation=https://github.com/luminos-os/luminos
After=network.target
# Power daemon should start before luminos-ai so the initial AC state is applied at boot.
Before=luminos-ai.service

[Service]
Type=simple
ExecStart=/usr/local/bin/luminos-power
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luminos-power
# Needs access to /sys/class/power_supply and /sys/class/hwmon.
RuntimeDirectory=luminos
RuntimeDirectoryMode=0755
# [CHANGE: claude-code | 2026-06-10] BUG: /run/luminos is SHARED by 5 daemons.
# Without Preserve, restarting any one of them wipes the dir and unlinks every
# other daemon's socket (happened 2026-06-08: power restart killed ai/sentinel/ram sockets).
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
