# config/hive-crash-watcher.path
[Unit]
Description=Watch for new coredumps

[Path]
PathChanged=/var/lib/systemd/coredump
Unit=hive-crash-watcher.service

[Install]
WantedBy=multi-user.target
